# launchd/com.hiretyler.jobbuddy.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hiretyler.jobbuddy.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/tylergeddes/projects/jobbuddy/server/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/tylergeddes/projects/jobbuddy</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/tylergeddes/projects/jobbuddy/logs/dashboard.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tylergeddes/projects/jobbuddy/logs/dashboard.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>America/Denver</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
